# pisces/forsense/dashboardinstall/systemd/vpn-check.service
[Unit]
Description=VPN checker

[Service]
ExecStart=/etc/monitor-scripts/vpn.sh
